# letsencrypt-nginx.timer
[Unit]
Description=Renew Let's Encrypt certificate (nginx)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=multi-user.target
